# ansible/templates/install_hdparm/hdparm.service.j2
[Unit]
Description=HDPARM HDD No Sleep

[Service]
Type=simple
ExecStart={{ retronas_root }}/scripts/hdparm.sh
TimeoutStopSec=10
Restart=never

[Install]
WantedBy=default.target
