# ansible/templates/install_hdparm/hdparm.timer.j2
[Unit]
Description=Run HDPARM HDD No Sleep

[Timer]
OnBootSec=3m
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target
