# ansible/templates/install_kermit/iksd.socket.j2
[Socket]
ListenStream = 1649
Accept=yes

[Install]
WantedBy = sockets.target