# ansible/templates/install_3ds_qr_codes/install_3ds_qr_codes.cron.j2
# RetroNAS automatically generated
# Every 15 minutes check for new 3DS titles and generate a QR code to match
*/15 * * * *   {{ retronas_user }}  {{ retronas_root }}/scripts/service_3ds_qr_codes.sh >/dev/nul 2>&1
