# ansible/templates/install_3ds_qr_codes/retronas_3ds_qr.cron.j2
# RetroNAS automatically generated
# Every 15 minutes check for new 3DS titles and generate a QR code to match
*/15 * * * *   {{ retronas_user }}  {{ retronas_root }}/scripts/3ds_qr.sh >/dev/null 2>&1
