# ansible/templates/install_romm_cifs/retronas-romm-dirs.service.j2
[Unit]
Description=retronas-romm-dirs

[Service]
Type=simple
Restart=no
ExecStart={{ retronas_root }}/scripts/retronas-romm-dirs.sh
TimeoutStartSec=0
RemainAfterExit=no
