# ansible/templates/install_linux-dexdrive/linux-dexdrive.service.j2
[Unit]
Description=Dex Drive service %i

[Service]
Type=forking
ExecStart=/usr/local/bin/dexattach -D \$DEVICE
TimeoutStopSec=10
Restart=no

[Install]
WantedBy=default.target
