# ansible/templates/install_hb-store-cdn/hb-store-cdn.cron.j2
0 */12 * * * {{ retronas_root }}/bin/{{ my_name }}/hb-store-cdn-check.sh