# ansible/templates/install_ps3netsrv/ps3netsrv-perm@.path.j2
[Unit]
Description=ps3netsrv permissions monitor %i


[Path]
PathChanged={{ retronas_path }}/ps3/ps3netsrv/%i/

[Install]
WantedBy=multi-user.target