# ansible/templates/install_ps3netsrv/ps3netsrv-perms.timer.j2
[Unit]
Description=ps3netsrv for PlayStation3 permissions fix

[Timer]
OnUnitActiveSec=20s
RandomizedDelaySec=15s
Persistent=true

[Install]
WantedBy=timers.target